# files/container/sshd.ini
[supervisord]
nodaemon=true

[program:sshd-keys]
command=/opt/generate_hostkey.sh
startsecs=0

[program:sshd]
command=/usr/sbin/sshd -D
